# scripts/crontab
0 9 * * 1-5 /bin/sh /app/scripts/cron.sh In
0 17 * * 1-5 /bin/sh /app/scripts/cron.sh Out
